# overlay/etc/cron.hourly/openvpn-profiles-delexpired
#!/bin/bash -e
# delete downloaded openvpn client profiles

su www-data -c /var/www/openvpn/bin/delexpired

